# stage2+pi-puck/04-swarmhack/files/pipuck_server.conf
[program:pipuck_server]
command=/home/pi/pipuck_server.sh
directory=/home/pi
autostart=true
startsecs=1
autorestart=true
stopasgroup=true
killasgroup=true
